# overlay/usr/lib/systemd/system/dhcpcd.service
# /usr/lib/systemd/system/dhcpcd.service
#
# systemd service unit for the dhcpcd DHCP client.

[Unit]
Description=dhcpcd DHCP Client
Wants=network.target
Before=network.target
After=network-pre.target systemd-udevd.service
Conflicts=NetworkManager.service
BindsTo=sys-subsystem-net-devices-%i.device
Documentation=man:dhcpcd(8)

[Service]
# Type=forking indicates that the service forks a child process and the parent exits.
# systemd will consider the service "started" when the parent exits.
Type=forking
#
# ExecStart specifies the command to execute to start the service.
# -b: Background dhcpcd
# -q: Quiet mode (suppress informative messages)
# -w: Wait for interface to be available before starting
# -B: Ensure dhcpcd is forked and running before exiting (useful with Type=forking)
# -M: Do not manage a hostname
# -C resolv,ntp: Do not touch resolv.conf or ntp.conf (if systemd-resolved/timesyncd manage them)
# -i %I: Operate on the interface specified by the instance (%I)
#
ExecStart=/usr/bin/dhcpcd -b -q -w -B -M -C resolv,ntp %I
ExecReload=/usr/bin/dhcpcd -n %I
ExecStop=/usr/bin/dhcpcd -k %I
#
# Restart policy: If the service fails, restart it.
Restart=on-failure
TimeoutStopSec=5
Nice=-5

[Install]
WantedBy=multi-user.target network.target
Alias=dhcpcd@.service
